# systemd/detection.service
[Unit]
Description=Home Smart Box Camera Detection (Hailo NPU + MQTT + MJPEG)
After=network.target docker.service mosquitto.service
Wants=docker.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/hailo-apps
ExecStart=/home/pi/hailo-apps/venv_hailo_apps/bin/python3 /home/pi/hailo-apps/detection_service.py --headless
Restart=on-failure
RestartSec=10
Environment=HOME=/home/pi

[Install]
WantedBy=multi-user.target
